# deployment/deb/systemd.service
[Unit]

Description=${description}
After=syslog.target network.target local-fs.target

[Service]

Type=notify
NotifyAccess=main
LimitNOFILE=${limit_no_file}
LimitCORE=infinity
Restart=on-failure
TimeoutStartSec=300s
RestartSec=5s
WatchdogSec=60s
KillMode=mixed

Environment="HOME=${app_root}"
Environment="CODE_LOADING_MODE=interactive"

${other_systemd_commands}

ExecStart=${app_name} foreground

[Install]
WantedBy=multi-user.target
